# projects/daily-email-automation/config/daily-email.service
[Unit]
Description=Daily Email Service to Hamza
Documentation=https://github.com/your-username/DevOps-Internship-at-QuizPlus/tree/main/projects/daily-email-automation
After=network.target network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=root
Group=root
ExecStart=/opt/daily_email/send_daily_email.sh
StandardOutput=journal
StandardError=journal

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

# Environment
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

[Install]
WantedBy=multi-user.target
